# profiles/debian-compat/overlay/etc/systemd/system/codex-firstboot.service
[Unit]
Description=CoLinux Compat First-Boot Setup
After=network.target network-online.target
Wants=network-online.target
Before=codex-disk-inventory.service codex-auto-update.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/first-boot
RemainAfterExit=yes
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
ConditionPathExists=!/persist/.first-boot-done

[Install]
WantedBy=multi-user.target
