# profiles/alpine/overlay/etc/init.d/codex-firstboot
#!/sbin/openrc-run
# /etc/init.d/codex-firstboot
# CodexOS Lite — First-boot setup service
#
# Runs on boot to initialize the appliance if this is the first time.
# Creates persistence structure, configures Codex CLI, and sets up cron jobs.
# Marks completion with /persist/.codexos-system/first-boot-done to skip later.

name="codex-firstboot"
description="CodexOS Lite first-boot initialization"

depend() {
    need localmount
    after hostname
    before codex-disk-inventory
    before codex-auto-update
}

FIRSTBOOT_DONE="/persist/.codexos-system/first-boot-done"
FIRSTBOOT_SCRIPT="/usr/local/bin/first-boot"
FIRSTBOOT_LOG="/persist/logs/first-boot.log"

start() {
    # Check if first boot has already been done
    if [ -f "${FIRSTBOOT_DONE}" ]; then
        einfo "First boot already completed — skipping"
        return 0
    fi

    ebegin "Running first-boot setup"

    # Ensure log directory exists
    mkdir -p /persist/logs

    # Run the first-boot script
    if [ -x "${FIRSTBOOT_SCRIPT}" ]; then
        "${FIRSTBOOT_SCRIPT}" >> "${FIRSTBOOT_LOG}" 2>&1
        local rc=$?

        if [ $rc -eq 0 ]; then
            eend 0 "First-boot setup completed successfully"
        else
            eend $rc "First-boot setup failed (see ${FIRSTBOOT_LOG})"
        fi
    else
        eend 1 "First-boot script not found: ${FIRSTBOOT_SCRIPT}"
    fi
}

stop() {
    # Nothing to stop
    return 0
}
