# profiles/alpine/overlay-desktop/etc/xdg/autostart/codex-desktop-autostart.desktop
[Desktop Entry]
Type=Application
Name=Codex Desktop
Comment=OpenAI Codex Desktop App
Exec=/opt/codex-desktop/codex-desktop %U
Icon=codex-desktop
Terminal=false
Categories=Development;IDE;
StartupNotify=true
X-GNOME-Autostart-enabled=true
