# profiles/alpine/overlay-desktop/etc/init.d/codex-network
#!/sbin/openrc-run
# CoLinux Desktop — codex-network init (OpenRC)
# Starts NetworkManager on boot for the desktop edition.

name="codex-network"
description="CoLinux Desktop network initialization (NetworkManager)"
command="/usr/sbin/NetworkManager"
command_background=true
pidfile="/run/NetworkManager.pid"

depend() {
    need localmount dbus
    after bootmisc
}
