# profiles/debian/overlay/etc/systemd/system/codex-update.service
[Unit]
Description=CoLinux Desktop Auto-Update Service
Documentation=https://github.com/gibavargas/colinux
After=network-online.target NetworkManager-wait-online.service
Wants=network-online.target
ConditionPathExists=/etc/codex-auto-update.conf

[Service]
Type=oneshot
ExecStart=/usr/local/bin/codex-system-autoupdate
User=root
Group=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codex-update
TimeoutStartSec=300
# Don't restart on failure — the timer handles periodic runs
Restart=no
# Resource limits
MemoryMax=512M
CPUQuota=50%

# Environment
Environment=CODEX_CONFIG_DIR=/etc
Environment=CODEX_LOG_DIR=/persist/logs
Environment=NODE_PATH=/usr/lib/node_modules

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target
