# profiles/debian/overlay/etc/systemd/system/codex-firstboot.service
[Unit]
Description=CodexOS Desktop First Boot Setup
Documentation=https://github.com/colinux/codexos
After=local-fs.target network.target NetworkManager.service
Wants=network-online.target
ConditionPathExists=!/persist/.codexos-system/first-boot-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/codex-firstboot-desktop
User=root
Group=root
StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=codex-firstboot
TimeoutStartSec=120
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
